# com.claude.discord-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.discord-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/claude-discord-bridge/.venv/bin/python3</string>
        <string>/path/to/claude-discord-bridge/bridge.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DISCORD_TOKEN</key>
        <string>YOUR_DISCORD_BOT_TOKEN</string>
        <key>CHANNEL_ID</key>
        <string>YOUR_CHANNEL_ID</string>
        <key>CLAUDE_TIMEOUT</key>
        <string>300</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/you</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/you/.claude-discord-bridge/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/you/.claude-discord-bridge/bridge.err.log</string>
</dict>
</plist>
